# sourcedigger-server.service
[Unit]
Description=SourceDigger Server

[Service]
Type=simple
ExecStart=/srv/sourcedigger/sourcedigger-server

[Install]
WantedBy=multi-user.target